# dot_Library/LaunchAgents/com.chezmoi.mac-dotfiles.maintenance.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.chezmoi.mac-dotfiles.maintenance</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>{{ .chezmoi.homeDir }}/.local/bin/mac-dotfiles-maintenance.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>{{ .data.maintenance.hour }}</integer>
        <key>Minute</key>
        <integer>{{ .data.maintenance.minute }}</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/mac-dotfiles/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/.local/state/mac-dotfiles/launchd.log</string>
</dict>
</plist>
